# crontab
MAILTO=francis@scraperwiki.com
@hourly ./tool/get.py
